# scripts/systemd/system/export-acmplus-csv.timer
[Unit]
Description=Vuelca datos de acmplus del dí­a anterior

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
