# scripts/systemd/system/export-pb200-csv.timer
[Unit]
Description=Vuelca datos de pb200 del dí­a anterior

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
