# scripts/systemd/system/reboot-computer.service
[Unit]
Description=Reinicia el sistema

[Service]
ExecStart=/usr/local/bin/reboot-computer